# etc/systemd/system/error-slack@.service
[Unit]
Description=Reports systemd timer error to Slack

[Service]
Type=oneshot
ExecStart=/usr/local/bin/systemd-error-slack %i